# debian/postinst
#!/bin/sh
update-rc.d snapperShadowCopy2 defaults
service snapperShadowCopy2 start