# debian/postrm
#!/bin/sh
update-rc.d snapperShadowCopy2 remove